# root/etc/s6/application/run
#!/usr/bin/env bash

exec /entrypoint.sh serve 